# deploy/launchd/com.iterwheel.voyager.loop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.iterwheel.voyager.loop</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>exec /Users/frank/Projects/voyager/scripts/loop_continue.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/frank/Projects/voyager</string>

  <key>StartInterval</key>
  <integer>900</integer>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/frank/Library/Logs/voyager/loop.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/frank/Library/Logs/voyager/loop.err.log</string>

  <key>Umask</key>
  <integer>63</integer>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/frank</string>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/frank/.local/state/fnm_multishells/default/bin:/opt/homebrew/bin</string>
    <key>VOYAGER_PROJECT_ROOT</key>
    <string>/Users/frank/Projects/voyager</string>
    <key>VOYAGER_ENV_FILE</key>
    <string>/Users/frank/.voyager/bridge.env</string>
    <key>VOYAGER_LOCK_FILE</key>
    <string>/tmp/voyager-loop.lock</string>
  </dict>
</dict>
</plist>
